# roles/systemd/files/rc.local
#!/bin/bash

/usr/local/sbin/unlock-disks.sh &

exit 0
